# 1.0/src/redis.sv.conf
[program:redis]
command=/usr/local/bin/start_redis.sh
autostart=true
autorestart=true
